# suspend-fix/resume-fix.service
[Unit]
Description=On Sleep enable
After=suspend.target

[Service]
ExecStart=/usr/local/bin/resume-mods
User=root

[Install]
WantedBy=suspend.target
