# clear-waydroid-desktop-entries/clear-waydroid-entries.service
[Unit]
Description=clear waydroid desktop entries

[Service]
ExecStart=/usr/local/bin/clear-waydroid-entries
Type=oneshot